# templates/default/akara.init.erb
#!/bin/sh

### BEGIN INIT INFO
# Provides:           akara-<%= @name %>
# Required-Start:
# Required-Stop:
# Default-Sart:       2 3 4 5
# Default-Stop:       0 1 6
# Short-Description:  Starts Akara (<%= @name %>) services
# Description:        Starts Akara (<%= @name %>) services
### END INIT INFO

APP_PATH=<%= @venv %>
DAEMON=${APP_PATH}/bin/akara
DAEMON_OPTS="-f ${APP_PATH}/akara.conf"
RUN_AS=<%= @user %>
PID_FILE=${APP_PATH}/logs/akara.pid

test -x $DAEMON || exit 1

set -e

case "$1" in
  start)
        echo -n "Starting Akara: "
        start-stop-daemon -d ${APP_PATH} -c ${RUN_AS} --start --background --pidfile ${PID_FILE} --startas ${DAEMON} -- ${DAEMON_OPTS} start
        ;;
  stop)
        echo -n "Stopping Akara: "
        start-stop-daemon --stop --pidfile ${PID_FILE}
        ;;
  status)
        ${DAEMON} ${DAEMON_OPTS} status
        ;;
  restart|force-reload)
        echo -n "Restarting Akara: "
        COMMAND="restart"
        BACKGROUND=""
        if [ ! -f "${PID_FILE}" ]; then
            COMMAND="start"
            BACKGROUND="--background"
        fi
        start-stop-daemon -d ${APP_PATH} -c ${RUN_AS} --start ${BACKGROuND} --pidfile ${PID_FILE} --startas ${DAEMON} -- ${DAEMON_OPTS} ${COMMAND}
        ;;
  *)
        echo "Usage: $0 {start|stop|status|restart|force-reload}" >&2
        exit 1
        ;;
esac

exit 0
